# lpcsnoop.service.in
[Unit]
Description=LPC Snoop Daemon
After=@SYSTEMD_AFTER_SERVICE@

[Service]
Restart=always
ExecStart=@bindir@/snoopd @SNOOPD_ARGS@

[Install]
WantedBy=@SYSTEMD_TARGET@
